# device/systemd/offline-ota-agent.service
[Unit]
Description=Offline OTA Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/offline-ota
ExecStart=/usr/bin/python3 -m agent.main status
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

