# device/systemd/offline-ota-dashboard.service
[Unit]
Description=Offline OTA Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/offline-ota
EnvironmentFile=-/etc/offline-ota/offline-ota.env
ExecStart=/usr/bin/python3 -m uvicorn server.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
